# .config/systemd/user/battery_notify.service
[Unit]
Description=Battery level notifier script

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery_notify.sh
